# systemd/strix-fc-agent.service
# SPDX-FileCopyrightText: 2026 Canonical, Ltd.
# SPDX-License-Identifier: GPL-2.0-only
#
# systemd service unit for the Strix FC Agent.
#
# The agent polls the Apollo Gateway for attachment state and reconciles
# the local strix_fc kernel module + iSCSI sessions accordingly.
#
# Install:
#   cp systemd/strix-fc-agent.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now strix-fc-agent
#
# Configuration is via environment variables.  Drop an override into
# /etc/systemd/system/strix-fc-agent.service.d/override.conf or use
# an EnvironmentFile.

[Unit]
Description=Strix FC Agent — gateway-driven FC+iSCSI reconciliation
Documentation=https://github.com/lunacy-systems/strix-fc
After=network-online.target iscsid.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/strix-fc-agent run
Restart=on-failure
RestartSec=5

# Required env vars — set via EnvironmentFile or systemd drop-in
EnvironmentFile=-/etc/default/strix-fc-agent

# Hardening
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
